# templates/etc/update-motd.d/30-banner.j2
#!/bin/bash
cat<< 'EOF'

{{ system_motd_banner }}

EOF
