# templates/etc/update-motd.d/static_msg.j2
#!/bin/bash
cat<< 'EOF'

{{ __system_motd_item.content | default(__system_motd_message) }}
EOF
