# templates/etc/update-motd.d/static_msg.j2
#!/bin/bash
cat<< 'EOF' | while IFS= read -r line; do echo -e "$line"; done

{{ __system_motd_item.content | default(__system_motd_message) }}
EOF
