# templates/etc/profile.d/shell_history.sh.j2
# Description: This file is used to configure the history settings for bash and zsh.
#
# WARNING: This file is managed by Ansible. Do not edit this file directly.

# Prevent doublesourcing
if [ -z "$BASHHISTORYSOURCED" ]; then
  BASHHISTORYSOURCED="Y"

  if [ -n "$BASH_VERSION" ]; then
    # configure history time format.
    # See bash(1) for more options
    HISTTIMEFORMAT="{{ system_bash_histtimeformat }}"

    # control history file name
    HISTFILE="{{ system_bash_histfile }}"

    HISTCONTROL="{{ system_bash_histcontrol }}"
    HISTIGNORE="{{ system_bash_histignore | join(':') }}"

    # append to the history file, don't overwrite it
  {% if system_bash_histappend %}
    shopt -s histappend
  {% endif %}
    # save all lines of a multiple-line command in the same history entry
    shopt -s cmdhist

    # for setting history length see HISTSIZE and HISTFILESIZE in bash(1)
    HISTSIZE={{ system_bash_histsize }}
    HISTFILESIZE={{ system_bash_histfilesize }}

    readonly HISTTIMEFORMAT
    readonly HISTFILE
    readonly HISTCONTROL
    readonly HISTSIZE
    readonly HISTFILESIZE
  {% if system_bash_prompt_command %}

    PROMPT_COMMAND='{{ system_bash_prompt_command }}'
  {% endif %}
  fi

  if [ -n "$ZSH_VERSION" ]; then
    HISTSIZE={{ system_bash_histsize }}
    SAVEHIST={{ system_bash_histfilesize }}
    HISTFILE=~/.zsh_history
    setopt APPEND_HISTORY
    setopt EXTENDED_HISTORY

    readonly HISTSIZE
    readonly SAVEHIST
    readonly HISTFILE
  {% if system_bash_prompt_command %}

    PROMPT_COMMAND='{{ system_bash_prompt_command }}'
  {% endif %}
  fi
fi

